# systemd/wgatectl.service
[Unit]
Description=wgatectl — LAN monitor & access-control daemon
After=network-online.target nss-lookup.target dnsmasq.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/wgatectl
Restart=on-failure
RestartSec=2s

User=root
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_DAC_READ_SEARCH CAP_CHOWN CAP_FOWNER
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/wgatectl
PrivateTmp=yes
MemoryMax=96M
LimitNOFILE=4096

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
